# .config/autostart/emacs-daemon.desktop
[Desktop Entry]
Name=Emacs Daemon
GenericName=Editor
Comment=Emacs Daemon
Exec=emacs --daemon
Terminal=false
Type=Application
